# config/com.tapo.onvif.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd template. The placeholders __BRIDGE_DIR__ and __HOME__ are
  substituted by `./install.sh launchd`, which also loads the agent.
  Manual install:
    sed -e "s|__BRIDGE_DIR__|$PWD|g" -e "s|__HOME__|$HOME|g" \
      config/com.tapo.onvif.plist.example \
      > ~/Library/LaunchAgents/com.tapo.onvif.plist
    launchctl load -w ~/Library/LaunchAgents/com.tapo.onvif.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tapo.onvif</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__BRIDGE_DIR__/src/run_bridge.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__BRIDGE_DIR__/src</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>__BRIDGE_DIR__/tmp/tapo-launchd.log</string>

  <key>StandardErrorPath</key>
  <string>__BRIDGE_DIR__/tmp/tapo-launchd.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>
</dict>
</plist>
