# config/com.tapo.bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd template. Copy to ~/Library/LaunchAgents/com.tapo.bridge.plist
  and edit the two paths to match your install location. Then:
    launchctl load -w ~/Library/LaunchAgents/com.tapo.bridge.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tapo.bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/YOU/tapo-bridge/src/run_bridge.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/YOU/tapo-bridge/src</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/tmp/tapo-launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/tapo-launchd.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/YOU</string>
  </dict>
</dict>
</plist>
